# templates/default/pgpool2-upstart.conf.erb
description "pgpool-II database connection manager"
author      "Pulselocker, Inc."

start on started networking
stop on shutdown

console owner

pre-start script
	if [ -d <%= node['pgpool2']['run_dir'] %> ]; then
		chmod 2775 <%= node['pgpool2']['run_dir'] %>
	else
		install -d -m 2775 -o <%= node['pgpool2']['user'] %> -g <%= node['pgpool2']['group'] %> <%= node['pgpool2']['run_dir'] %>
	fi
end script

exec sudo -u <%= node['pgpool2']['user'] %> <%= node['pgpool2']['prefix_dir'] %>/bin/pgpool -n -f <%= node['pgpool2']['config_dir'] %>/pgpool.conf > <%= node['pgpool2']['log_dir'] %>/pgpool.log 2>&1

respawn
